# files/etc/services.d/encryption/finish
#!/usr/bin/execlineb -S0

/stopscripts/shutdown.sh

s6-svscanctl -t /var/run/s6/services