# files/etc/services.d/encryption/run
#!/usr/bin/with-contenv bash
/startscripts/encryption.sh